# attitude_adjustment/feeds/packages/ipv6/ptrtd/files/ptrtd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=60

PREFIX="3ffe:abcd:1234:9876::"
LEN=64

start() {
	service_start /usr/sbin/ptrtd -p $PREFIX -l $LEN
}

stop() {
	service_stop /usr/sbin/ptrtd
}
