# attitude_adjustment/feeds/packages/.svn/pristine/6b/6bc637908120b100e32014ec60d51eeca8c3dae3.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2011 OpenWrt.org

START=99

SERVICE_DAEMONIZE=1

start()
{
	/usr/sbin/chaosvpn_config.sh
	[ $? = 0 ] || {
		logger "not starting chaosvpn, check your config"
		exit 1
	}
	service_start /usr/sbin/chaosvpn -f -c /etc/tinc/chaosvpn.conf
}

stop()
{
	service_stop /usr/sbin/chaosvpn
}
