# attitude_adjustment/feeds/packages/net/ptpd/files/ptpd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=50

start() {
	service_start /usr/sbin/ptpd
}

stop() {
	service_stop /usr/sbin/ptpd
}
