# attitude_adjustment/target/linux/armadaxp/base-files/etc/init.d/mamba
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org 

START=09
start() {
	echo 33 > /sys/class/gpio/unexport
	sleep 1
	echo 40 > /sys/class/gpio/export
	echo 24 > /sys/class/gpio/export
	sleep 1
	/usr/sbin/app_reset start /usr/sbin/leds_monitor

	/usr/sbin/wan_monitor&
}
