# attitude_adjustment/feeds/packages/.svn/pristine/90/9024d1c08d520d1678e3dfa5d4aba0d28c91a5d1.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2012 OpenWrt.org

START=50

SERVICE_USE_PID=1
SERVICE_PID_FILE=/var/run/pepper.pid

start() {
	. /lib/functions/network.sh

	local ipaddr
	network_get_ipaddr ipaddr lan && \
		service_start /usr/sbin/pepper --dns1="$ipaddr"
}

stop() {
	service_stop /usr/sbin/pepper
}
