# attitude_adjustment/feeds/xwrt/webif-asterisk/etc/crontab/root
57 23 * * * /www/cgi-bin/webif/send_call_log.sh > /dev/null 2>&1
0-59 * * * * /www/cgi-bin/webif/pid_send_all.sh > /dev/null 2>&1



