# attitude_adjustment/feeds/packages/.svn/pristine/8c/8c62043dd82b26d19803c680ed257a3fe4adec01.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=50

start() {
	user_exists scanlogd 51 || user_add scanlogd 51
	group_exists scanlogd 51 || group_add scanlogd
	[ -d  /var/run/scanlogd ] || {
		mkdir -m 0755 -p /var/run/scanlogd
		chmod 0755 /var/run/scanlogd
		chown scanlogd:scanlogd /var/run/scanlogd
	}
	service_start /usr/sbin/scanlogd
}

stop() {
	service_stop /usr/sbin/scanlogd
}
