# attitude_adjustment/feeds/packages/.svn/pristine/71/71cdc6cba271dd880f238ef50b6db94d29e7b7f6.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

start() {
	$BIN -b $PORT &
}

stop() {
	killall php
}

