# attitude_adjustment/feeds/packages/utils/psplash/files/psplash.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=05

start() {
    psplash &
}

stop() {
    killall psplash
}
