# attitude_adjustment/feeds/packages/ipv6/dibbler/files/dibbler-relay.init
#!/bin/sh /etc/rc.common

START=90
STOP=90

start() {
	mkdir -p /var/lib/dibbler
	dibbler-relay start
}

stop() {
	dibbler-relay stop
}
