# attitude_adjustment/feeds/packages/libs/howl/files/autoipd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

# XXX: pid-file is /var/run/autoipd-<INTERFACE>.pid
SERVICE_USE_PID=

start() {
	service_start /usr/sbin/autoipd
}
stop() {
	service_stop /usr/sbin/autoipd
}

