# attitude_adjustment/feeds/packages/.svn/pristine/33/33636dbab93d0f47928d0f922d9c7d4852c26923.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2012 OpenWrt.org

START=90

start() {
	freecwmpd 2>/dev/null >/dev/null &
}

stop() {
	killall -9 freecwmpd 2>/dev/null >/dev/null
}
