# attitude_adjustment/feeds/packages/.svn/pristine/32/321cdc11bb60acf0aaca42a11cef26ab380adbfb.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=94

SERVICE_USE_PID=1

start() {
	mkdir -m 0755 -p /var/lib/scmpc
	service_start /usr/bin/scmpc
}

stop() {
	service_stop /usr/bin/scmpc
}
