# attitude_adjustment/feeds/packages/.svn/pristine/82/8228783f801f5acaee41a2e914925e9d96e3425e.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

start() {
	mkdir -m 0755 -p /var/lib/bitlbee
	service_start /usr/sbin/bitlbee -D
}

stop() {
	service_stop /usr/sbin/bitlbee
}
