# attitude_adjustment/feeds/packages/net/arpwatch/files/arpwatch.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=90

start() {
	service_start /usr/sbin/arpwatch -f /etc/arpwatch/arp.dat -i br-lan
}

stop() {
	service_stop /usr/sbin/arpwatch
}
