# attitude_adjustment/feeds/packages/.svn/pristine/a4/a4fad4e9d373c9e6be9009ccd4c5717f91851b1c.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

SERVICE_USE_PID=1

PROG="uacctd"
DEFAULT="/etc/default/$PROG"
COMMAND="/usr/sbin/$PROG"

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	service_start $COMMAND $OPTIONS
}

stop() {
	service_stop $COMMAND
}
