# attitude_adjustment/feeds/packages/.svn/pristine/36/36f3a9e69f50af2dd073d530e33d9613125eaa1b.svn-base
#!/bin/sh /etc/rc.common 

START=50

start() {
	LCDd -c /etc/LCDd.conf
}

stop() {
	killall LCDd
}
