# attitude_adjustment/feeds/packages/.svn/pristine/86/867ef941e21f03f0fc05dc2009fab280bdbea607.svn-base
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

PROG="ulogd"
COMMAND="/usr/sbin/$PROG"
DEFAULT="/etc/default/$PROG"

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	service_start $COMMAND $OPTIONS
}

stop() {
	service_stop $COMMAND
}
