# attitude_adjustment/feeds/packages/utils/sysstat/files/sysstat.init
#!/bin/sh /etc/rc.common

# The init file runs the sadc command in order to write
# the "LINUX RESTART" mark to the daily data file

START=99

start() {
	mkdir -p /var/log/sysstat
	/usr/lib/sysstat/sadc -S DISK -F -L -
}

