# attitude_adjustment/feeds/packages/.svn/pristine/38/38743a984aaa2013296aa063c58ab183734df663.svn-base
#!/bin/sh /etc/rc.common
START=99
EXTRA_COMMANDS="single"

start () {
   /usr/bin/multiwan agent &
}

stop () {
  sh /usr/bin/multiwan stop
}

restart () {
  /usr/bin/multiwan restart &
}

single () {
  /usr/bin/multiwan single &
}
  
