# attitude_adjustment/feeds/packages/.svn/pristine/59/5906668816f5b1825746c7d98d28964305beab39.svn-base
#!/bin/sh /etc/rc.common. 

START=50

start() {
	lcdvc -c /etc/lcdvc.conf
}

stop() {
	killall lcdvc
}
