# attitude_adjustment/feeds/packages/.svn/pristine/76/76846c2590a457f4b57512766dc99b9a7603f856.svn-base
#!/bin/sh /etc/rc.common

START=98
start() {
	# make /tmp/.bitchxrc as many servers ban root
	if [ ! -e /tmp/.bitchxrc ]
	then
		echo "IRCUSER OpenWrt" >/tmp/.bitchxrc
		echo "IRCNAME OpenWrt User!" >>/tmp/.bitchxrc
	fi
}
