# attitude_adjustment/feeds/packages/.svn/pristine/2c/2cceb086027c4c46afc4005f54e641700c8a4e5c.svn-base
#!/bin/sh /etc/rc.common
#Copyright (C) 2007 RoadRunner.cx
START=85

start() {
    [ -f /etc/imsnif.conf ] && /usr/bin/imsnif || {
	cat <<EOF
You have not created configuration file.
Store your configuration as /etc/imsnif.conf
EOF
    }
}

stop() {
    [ `ps ax|grep -v grep|grep -c imsnif` != "0" ] && killall imsnif
}