# attitude_adjustment/feeds/packages/.svn/pristine/3e/3e96576e1c50b003aa5450568f0a43dcb04a2564.svn-base
#!/bin/sh /etc/rc.common. 

START=50

start() {
	lcdproc -c /etc/lcdproc.conf
}

stop() {
	killall lcdproc
}
